# deployment/picoclaw.service
[Unit]
Description=PicoClaw AI Agent
After=network.target
After=ollama.service

[Service]
Type=simple
# TODO: Change 'racer' to your username
User=racer
# TODO: Change '/home/racer/picoclaw' to your installation path
WorkingDirectory=/home/racer/picoclaw
ExecStart=/home/racer/picoclaw/picoclaw
Restart=always
RestartSec=10
EnvironmentFile=/home/racer/picoclaw/.env

[Install]
WantedBy=multi-user.target
